# openclaw/skills/flight_search/flight-search.timer
[Unit]
Description=Daily OpenClaw flight search

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=flight-search.service

[Install]
WantedBy=timers.target
