# openclaw/skills/job_search_suhani/job-search-suhani.timer
[Unit]
Description=Daily retail job search for Suhani (07:30)

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true
Unit=job-search-suhani.service

[Install]
WantedBy=timers.target
